# weavy-manager.conf
[program:weavy-manager]
command=node /home/ubuntu/auto-credit-manager/backend/server.js
directory=/home/ubuntu/auto-credit-manager/backend
user=ubuntu
autostart=true
autorestart=true
stderr_logfile=/var/log/weavy-manager.err.log
stdout_logfile=/var/log/weavy-manager.out.log
environment=NODE_ENV=production